# docker/etc/s6/gitea/run
#!/bin/bash
[[ -f ./default_env ]] && source ./default_env
[[ -f ./setup ]] && source ./setup

pushd "${RUN_DIR}" > /dev/null
    exec su-exec $USER /app/gitea/gitea web
popd
