# docker/root/etc/s6/gitea/finish
#!/usr/bin/execlineb -S0

s6-svscanctl -t /var/run/s6/services
